# etc/systemd/system/docker.service.d/http-proxy.conf
[Service]    
Environment=HTTP_PROXY="http://127.0.0.1:7182/" 
Environment=HTTPS_PROXY="http://127.0.0.1:7182/" 
Environment=NO_PROXY="docker.nihaomai.com,localhost,127.0.0.1" 

